# 98-hdmi-monitor.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/usr/local/bin/on-monitor-change.sh"
